# java8090-742.ini
[program:meeting-server]
command=/usr/bin/java -Xms256m -Xmx512m -jar /home/wby-u24/trae/meeting-transcribe-1.0.0.jar --spring.config.location=/home/wby-u24/trae/application.yml
directory=/home/wby-u24/trae
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=100
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=20
stdout_logfile=/home/wby-u24/trae/logs/meeting-server.log
stderr_logfile=/home/wby-u24/trae/logs/meeting-server.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=20
stderr_capture_maxbytes=1MB
stderr_events_enabled=false
stopasgroup=false
killasgroup=false
