# supervisord-meeting.ini
[program:meeting-server]
directory=/opt/wechat_meeting
command=/usr/bin/java -Xms256m -Xmx512m -jar meeting-transcribe-1.0.0.jar --spring.config.location=/opt/wechat_meeting/application.yml
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0
stopsignal=TERM
stopwaitsecs=60
user=root
redirect_stderr=true
stdout_logfile=/var/log/wechat_meeting/meeting-server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=JAVA_HOME="/usr/lib/jvm/java-1.8.0-openjdk"
