# supervisord-tencent.ini
[program:meeting-server]
command=/usr/bin/java -Xms256m -Xmx512m -jar /opt/wechat_meeting/meeting-transcribe-1.0.0.jar --spring.config.location=/opt/wechat_meeting/application.yml
directory=/opt/wechat_meeting
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=1
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=20
stdout_logfile=/var/log/wechat_meeting/meeting-server.log
stderr_logfile=/var/log/wechat_meeting/meeting-server.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=20